# sysperf.service

[Unit]
Description= SysPerf backend

[Service]
User= root
ExecStart=/bin/SysModule
Restart = always

[Install]
WantedBy=multi-user.target

